# etc/systemd/system/music-mount.service
[Unit]
Description=Mount USB Music Drive
After=dev-music.device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/music-mount start

# [Install]
# WantedBy=multi-user.target
