# etc/systemd/system/music-umount.service
[Unit]
Description=Stop mpd when music USB is removed

[Service]
Type=oneshot
ExecStart=/usr/local/bin/music-umount

